# etc/udev.rules
ACTION=="add" \
, ENV{DEVTYPE}=="disk" \
, KERNEL=="sd*" \
, RUN+="/bin/bash -c 'sudo -i -u pi DISPLAY=:0 wmctrl -c ESCAPE; sudo -i -u pi DISPLAY=:0 chromium-browser --app=http://localhost:8000/0$(ls /dev/sd{a,b,c,d} | wc -l)/; mpg123 --verbose /home/pi/Documents/escape-2023/public/scenes/0$(ls /dev/sd{a,b,c,d} | wc -l)/audio.mp3'"
